# git_home/net-cgi.git/files/dns_hijack.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=98

start() {
	hijack_process="$(/bin/config get hijack_process)"

	if [ "$hijack_process" = "1" ]; then
		/bin/config set dns_hijack="1"
		/bin/config commit
		/usr/sbin/dns-hijack
	fi
}

stop() {
	echo "\n"
}

